# files/services/codecserver/run
#!/usr/bin/execlineb -P
/usr/bin/codecserver
